# lcampuzano/scripts/stop-dos.service
[Unit]
Description=Mi servicio super-útil

[Service]
ExecStart=/home/lcampuzano/scripts/dos-stop.sh

[Install]
WantedBy=multi-user.target
